# deploy/systemd-user-optional/xcelium-mcp-culler.service
# Companion to xcelium-mcp-culler.timer — see that file and
# design.md §7.3 for the (optional, linger-gated) systemd --user promotion path.

[Unit]
Description=xcelium-mcp idle worker culler (one-shot, triggered by the .timer)

[Service]
Type=oneshot
# 2026-07-07: see xcelium-mcp-supervisor.service — same root-owned-bin-dir reason.
ExecStart=/opt/mcp-env/bin/python3 -m xcelium_mcp.idle_culler
